# templates/default/init.d.bash.erb
#!/bin/bash
#
# This file was generated by chef for <%= node['fqdn'] %>
# DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN BY CHEF CLIENT
#
### BEGIN INIT INFO
# Provides: <%= @service_name %>
# Required-Start: $network $remote_fs $syslog $local_fs
# Required-Stop: $network $remote_fs $syslog $local_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: starts <%= @service_name %> by eye
### END INIT INFO


EYE_BIN=<%= @eye_bin %>
CONFIG_DIR=<%= @config_dir %>
RUN_USER=<%= @user || 'root' %>
LOG_FILE=<%= @log_file || '/var/log/eye.log' %>

source /etc/eye/functions

case "$1" in
  start)
    load
    load_all
    ;;
  stop)
    quit
    ;;
  stop-all)
    quit true
    ;;
  restart)
    quit
    load
    load_all
    ;;
  reload)
    load_all
    ;;
  status)
    status
    ;;
  full-status)
    status '' true
    ;;
  watch)
    execute watch all
    ;;
  explain)
    check_all true true
    ;;
  break)
    execute break all
    ;;
  check)
    check_all false true
    ;;
  *)
    echo "Usage: $0 {start|stop|stop-all|restart|status|full-status|check|explain|watch|break}"
    exit 1
    ;;
esac

