# templates/default/init.d.bash.erb
#!/bin/bash
#
# This file was generated by chef for <%= @node[:fqdn] %>
#
### BEGIN INIT INFO
# Provides: <%= @service_name %>
# Required-Start: $network $remote_fs $syslog $local_fs
# Required-Stop: $network $remote_fs $syslog $local_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: starts <%= @service_name %> by eye
### END INIT INFO


EYE_BIN=<%= @node['chef_eye']['bin'] %>
CONFIG_DIR=<%= @config_dir %>
RUN_USER=<%= @user || 'root' %>
LOG_FILE=<%= @log_file || '/var/log/eye.log' %>
MASK=$2

source /etc/eye/functions

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    stop
    start
    ;;
  reload)
    reload $MASK
    ;;
  status)
    status
    ;;
  full-status)
    status true
   ;;
  watch)
    execute watch $MASK
   ;;
  explain)
   if [ ! $MASK ]; then
      echo -e "${RED}CONF is required!${CLEAR}"
   else
      explain $MASK
      echo -e $OUTPUT
   fi
   ;;
  break)
    if [ ! $MASK ]; then
      echo -e "${RED}MASK is required!${CLEAR}"
    else
      execute break $MASK
    fi
   ;;
  check)
    check_all
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status|full-status|check|explain|watch|break} [MASK|CONF]"
    exit 1
    ;;
esac

